# deploy/organizator3000.service
[Unit]
Description=Organizador Pessoal (task organizer)
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/organizator3000/server
ExecStart=/usr/bin/node dist/index.js
Environment=NODE_ENV=production
Environment=PORT=4000
# Config extra de controle de acesso (Cloudflare Access) fica aqui, fora do git.
# O "-" faz o serviço subir normalmente mesmo se o arquivo não existir.
EnvironmentFile=-%h/organizator3000/.env
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
